# nginx.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=www-data:www-data

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/php-fpm-supervisor.log
stderr_logfile=/var/log/supervisor/php-fpm-supervisor.err

[program:nginx]
command=nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-supervisor.log
stderr_logfile=/var/log/supervisor/nginx-supervisor.err
